# services/minikube.service
[Unit]
Description=MiniKube dashboard
After=network.target

[Service]
Type=simple
User=minikube
ExecStart=/usr/bin/minikube start
ExecStop=/usr/bin/minikube stop
KillMode=process
Restart=on-failure
RestartSec=42s

[Install]
WantedBy=multi-user.target